# root/etc/services.d/tvapp2/run
#!/usr/bin/with-contenv sh

# #
#   install and startup for tvapp2
# #

cp -r ${DIR_BUILD} ${DIR_RUN}
rm -rf ${DIR_BUILD}/*
cd ${DIR_RUN}
npm start
