# root/etc/services.d/tvapp2/run
#!/usr/bin/with-contenv sh

# #
#   install and startup for tvapp2
# #

cp -r ${INSTALL_DIR} ${WORKING_DIR}
rm -rf ${INSTALL_DIR}/*
cd ${WORKING_DIR}
npm start
